# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/lib/systemd/system/button.service
[Unit]
Description=Button Reboot Service
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/buttonservice.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
