# foxbuntu/sysdrv/out/rootfs_uclibc_rv1106/lib/systemd/system/settime.service
[Unit]
Description=Run settime script at startup
Before=meshtasticd.service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/settime.sh

[Install]
WantedBy=default.target
